# api/supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm]
command = php-fpm

[program:splunkd]
command = /opt/splunkforwarder/bin/splunk start --nodaemon

